# VM1/etc/udev/rules.d/70-snap.snap-store.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_snap-store_snap-store"
# opengl
KERNEL=="nvmap", TAG+="snap_snap-store_snap-store"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_snap-store_snap-store"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_snap-store_snap-store"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_snap-store_snap-store"
# opengl
KERNEL=="vchiq", TAG+="snap_snap-store_snap-store"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_snap-store_snap-store"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_snap-store_snap-store"
TAG=="snap_snap-store_snap-store", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_snap-store_snap-store $devpath $major:$minor"
# opengl
KERNEL=="nvhost-*", TAG+="snap_snap-store_ubuntu-software"
# opengl
KERNEL=="nvmap", TAG+="snap_snap-store_ubuntu-software"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_snap-store_ubuntu-software"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_snap-store_ubuntu-software"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_snap-store_ubuntu-software"
# opengl
KERNEL=="vchiq", TAG+="snap_snap-store_ubuntu-software"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_snap-store_ubuntu-software"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_snap-store_ubuntu-software"
TAG=="snap_snap-store_ubuntu-software", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_snap-store_ubuntu-software $devpath $major:$minor"
# opengl
KERNEL=="nvhost-*", TAG+="snap_snap-store_ubuntu-software-local-file"
# opengl
KERNEL=="nvmap", TAG+="snap_snap-store_ubuntu-software-local-file"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_snap-store_ubuntu-software-local-file"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_snap-store_ubuntu-software-local-file"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_snap-store_ubuntu-software-local-file"
# opengl
KERNEL=="vchiq", TAG+="snap_snap-store_ubuntu-software-local-file"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_snap-store_ubuntu-software-local-file"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_snap-store_ubuntu-software-local-file"
TAG=="snap_snap-store_ubuntu-software-local-file", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_snap-store_ubuntu-software-local-file $devpath $major:$minor"
